# units/platformd-trustd.service.in
[Unit]
Description=platformd platform-authentication state authority
Documentation=man:platformd-trustd.service(8)

[Service]
Type=notify
ExecStart=@LIBEXECDIR@/platformd-trustd
RuntimeDirectory=platformd-trustd
RuntimeDirectoryMode=0755
StateDirectory=platformd-trustd
StateDirectoryMode=0700

# Observability only: the daemon reads system state and serves it read-only, so
# it can run with almost nothing. Tighten hard.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
